# miaosha.supervisor.ini
[program:miaosha]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/php/bin/php /data/wwwroot/artisan queue:work --queue=generate-order
autostart=true
autorestart=true
user=root
numprocs=4
redirect_stderr=true
stdout_logfile=/data/wwwlogs/miaosha.supervisor.log